# config/crontab/crontab
0	*	*	*	*	root	bash /home/www/<MY_APP_NAME>.test/www/bin/console.sh dump >/dev/null 2>&1